# etc/systemd/system/getty@tty1.service.d/10-thyriaen-autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --skip-login --login-options "-f thyriaen" --noreset --noclear --noissue %I $TERM
