# etc/systemd/system/getty@tty1.service.d/10-thyriaen-autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin thyriaen --noissue %I $TERM
