# cfg/systemd/user/thyachieve.timer
[Unit]
Description=Run thyachieve daily at 20:00

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
