# cfg/systemd/user/thyachieve.service
[Unit]
Description=Update Eww achievement progression

[Service]
Type=oneshot
Environment=XDG_SESSION_TYPE=wayland
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin
ExecStart=%h/.local/bin/thyachieve
