# cfg/systemd/user/thyachieve-work.service
[Unit]
Description=Update Eww achievement work status

[Service]
Type=oneshot
Environment=XDG_SESSION_TYPE=wayland
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin
ExecStart=%h/.local/bin/thyachieve-work
