# etc/systemd/system/sddm.service.d/10-cursor.conf
[Service]
Environment=XCURSOR_THEME=material_light_cursors
Environment=XCURSOR_SIZE=24
